# flightcontrol/templates/apiservice.service
[Unit]
Description=Smartbox rkt api-service
[Service]
Type=simple
ExecStart=/bin/sh -c 'rkt api-service'
Restart=always
KillMode=mixed