# build/smartbox.service
[Unit]
Description=Smartbox flightcontrol
[Service]
Type=simple
ExecStartPre=/bin/sh -c '/bin/mkdir -p $SMARTBOX_HOME/com'
ExecStart=/bin/sh -c 'rkt run --insecure-options=image --stage1-name=coreos.com/rkt/stage1-fly:1.4.0 --net=none --volume unitfiles,kind=host,source="/etc/systemd/system" --volume com,kind=host,source="$SMARTBOX_HOME/com" --volume dbussocketlocation,kind=host,source="/var/run/dbus" --volume root,kind=host,source="$SMARTBOX_HOME/flightcontrol" $SMARTBOX_HOME/flightcontrol.aci'
Restart=yes
KillMode=mixed

